# blocks/GM901/etc-init-GM901.conf
description  "Device Service for Model GM901"

start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn

env NODE_ENV=production
exec /usr/local/bin/node /root/gm901/main.js >> /var/log/gm901.log 2>&1
