# blocks/Weiran/etc-init-Weiran.conf
description  "Weiran Web Service"

start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn

env NODE_ENV=production
exec /usr/local/bin/node /root/weiran/weiran.js >> /root/weiran/weiran.log 2>&1
